# bitrix/supervisor/push/supervisord.ini
[supervisord]
nodaemon=true

[program:push-server-multi-reset]
command = /etc/init.d/push-server-multi reset

[program:push-server-multi-start]
command = /etc/init.d/push-server-multi start
autorestart=true             

[program:redis-server]
command = /usr/bin/redis-server /etc/redis.conf --daemonize no
autorestart=true             

